# services/LaunchTransmission.service
[Unit]
Description=Start Transmission minimized

[Service]
Type=oneshot
ExecStart=/usr/bin/transmission-gtk -m

[Install]
WantedBy=graphical.target